# etc/systemd/system/rtkit-daemon.service.d/log.conf
[Service]
# Limit logging to "warning"
LogLevelMax=4
