# etc/systemd/system/home-narcizo-mnt-castor.automount
[Unit]
Description=sshfs automount on /home/narcizo/mnt/castor

[Automount]
Where=/home/narcizo/mnt/castor
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
