# etc/systemd/system/media-Pictures.automount
[Unit]
Description=CIFS automount (Pictures on elefante)

[Automount]
Where=/media/Pictures
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
