# etc/systemd/system/media-Videos.automount
[Unit]
Description=CIFS automount (Videos on elefante)

[Automount]
Where=/media/Videos
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
