# etc/systemd/system/media-Music.automount
[Unit]
Description=CIFS automount (Music on elefante)

[Automount]
Where=/media/Music
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
