# etc/systemd/system/mnt-backup.automount
[Unit]
Description=CIFS automount (Backup-linux on elefante)

[Automount]
Where=/mnt/backup
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
